# bin/v2ray-web.service
[Unit]
Description=V2ray Web Service
After=network-online.target
Wants=network-online.target
 
[Service]
Type=simple
ExecStart=/opt/v2ray-web/v2ray-web web
KillMode=control-group
Restart=on-failure
 
[Install]
WantedBy=multi-user.target